# rootfs/etc/s6-overlay/s6-rc.d/20-webproxy/up
#!/bin/sh
exec /etc/s6-overlay/scripts/20-webproxy
